# docker-mass-import-CSV/crontab
*/5 * * * * python /app/massActualCSVImport.py  # Runs every 5 minutes
